# resources/setup/mac/claude-setup/configs/com.claude.vault-autosync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude.vault-autosync</string>
  <key>ProgramArguments</key>
  <array>
    <string>__CLAUDE_DIR__/vault-autosync.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__CLAUDE_DIR__/vault-autosync.log</string>
  <key>StandardErrorPath</key>
  <string>__CLAUDE_DIR__/vault-autosync.log</string>
</dict>
</plist>
